# systemd/fedbox.socket.in
[Unit]
Description=FedBOX

[Socket]
ListenStream=LISTEN_HOST:LISTEN_PORT

[Install]
WantedBy=sockets.target

